# scraper/news-scraper.timer
[Unit]
Description=Run News Scraper every 5 hours

[Timer]
# Run every 5 hours
OnBootSec=5min
OnUnitActiveSec=5h
# Run immediately after boot if a scheduled run was missed
Persistent=true
# Random delay up to 5 minutes to avoid exact timing
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
